# install_src/khipu.service
[Unit]
Description="Khipu backend"

[Service]
ExecPaths="/opt/khipu/"